# infra/docker-backup.service
[Unit]
Description=Backup all Docker containers

[Service]
Type=oneshot
ExecStart=/docker/infra/main_backup.py full
User=tekjoey
Group=tekjoey

